# parts/test1/assets/misc/usb/p2/etc/rc.local
#!/bin/bash
[ -f "${0}.yac" ] && source ${0}.yac
exit 0
